# usb.service
[Service]
ExecStart=/home/ubuntu/Music/usbDetectServ

